# sh/crontab

*/30 10-21 * * *  $HOME/collector/collector.sh